# systemd/app.service
[Unit]
Description=WatchMyBirds Main Application
After=network-online.target go2rtc.service wmb-first-boot.service
Wants=network-online.target go2rtc.service wmb-first-boot.service
[Service]
User=watchmybirds
# Standard Systemd State/Logs management
StateDirectory=watchmybirds
LogsDirectory=app
WorkingDirectory=/opt/app
EnvironmentFile=-/etc/app/app.env
Environment=HOME=/opt/app/data
Environment=PORT=8050
Environment=OUTPUT_DIR=/opt/app/data/output
Environment=INGEST_DIR=/opt/app/data/ingest
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
Environment=MODEL_BASE_PATH=/opt/app/data/models
Environment=GO2RTC_CONFIG_PATH=/opt/app/data/output/go2rtc.yaml
# Force WiFi unblock on every start (runs as root via +). Sleep to avoid NM race.
ExecStartPre=+/bin/sleep 2
ExecStartPre=+/usr/sbin/rfkill unblock wifi
ExecStartPre=+/bin/sh -c 'if [ -f /etc/app/app.env ]; then chown root:root /etc/app/app.env && chmod 600 /etc/app/app.env; fi'
# Ensure app user can update/read go2rtc config after setup/import operations.
ExecStartPre=+/bin/sh -c 'if [ -f /opt/app/data/output/go2rtc.yaml ]; then chown watchmybirds:watchmybirds /opt/app/data/output/go2rtc.yaml && chmod 640 /opt/app/data/output/go2rtc.yaml; fi'
ExecStart=/opt/app/.venv/bin/python /opt/app/main.py

Restart=always
RestartSec=5

# Security Hardening (strict sandboxing)
ProtectSystem=strict
ReadWritePaths=/opt/app/data /var/log/app
PrivateTmp=true
# Power actions use logind/polkit; keep NNP enabled to block SUID escalation.
NoNewPrivileges=true
ProtectHome=yes
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
# Hardening: These flags are safe with logind/polkit power management (no sudo needed)
RestrictNamespaces=true
MemoryDenyWriteExecute=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
# SystemCallArchitectures=native  # Commented: Implicitly sets NoNewPrivileges, breaking sudo

[Install]
WantedBy=multi-user.target
